# root/etc/profile.d/transmission.sh
export PATH=$PATH:/opt/transmission/bin
