# root/etc/s6-overlay/s6-rc.d/transmission-init-config/run
#!/command/with-contenv sh

if [ ! -f /config/settings.json -o -s /config/settings.json ]; then
  cp /opt/defaults/transmission/settings.json /config/settings.json
fi