# root/etc/s6-overlay/s6-rc.d/transmission-init-config/run
#!/usr/bin/with-contenv sh

if [ ! -f /config/settings.json ]; then
  cp /opt/defaults/transmission/settings.json /config/settings.json
fi