# root/etc/s6-overlay/s6-rc.d/transmission-configure-anti-brute-force/up
/etc/s6-overlay/s6-rc.d/transmission-configure-anti-brute-force/run