# root/etc/s6-overlay/s6-rc.d/transmission-configure-blocklist/up
/etc/s6-overlay/s6-rc.d/transmission-configure-blocklist/run
